# install-stormos-xfce/airootfs/etc/skel/.config/autostart/welcome.desktop
[Desktop Entry]
Name=Storm Welcome
Comment=Get started with StormOS — Arch made calm
Exec=storm_welcome
Icon=system-software-install
Terminal=false
Type=Application
Categories=System
